# timeServer/supervisord/timeserver.conf
[program:timeserver]
directory=/root
command=node main.js

